# services/mysql/run
#!/bin/sh

# run service
exec /usr/sbin/mysqld "$MYSQL_OPTS"
