# rootfs/etc/services.d/mysql/run
#!/usr/bin/with-contenv sh

exec \
    /usr/sbin/mysqld $MYSQL_OPTS
